# systemd/docker-prune.timer
[Unit]
Description=Weekly Docker image prune

[Timer]
OnCalendar=Sun 04:30
Persistent=true

[Install]
WantedBy=timers.target
