# systemd/stack-backup.service
[Unit]
Description=Backup docker stack stack
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/stack-backup.sh
