# .codex-tmp/stack-autostart.service
[Unit]
Description=Start docker compose stacks in a controlled order after boot
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/stack-autostart.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
